# vjp13/.config/autostart/nocaps.desktop
[Desktop Entry]
Encoding=UTF-8
Name=nocaps
Comment=Remap caps lock as control.
Exec=/usr/bin/setxkbmap -option ctrl:nocaps
Terminal=true
