# vjp13/.config/autostart/synclient.desktop
[Desktop Entry]
Encoding=UTF-8
Name=synclient
Comment=Disable tap2click
Exec=/usr/bin/synclient MaxTapTime=0
Terminal=true
